# etc/init/whoiscache.conf

#
# WhoisCache Startup
#

description "WhoisCache"
author "Matthias Hannig <matthias@hannig.cc>"

start on runlevel [2345]
stop on runlevel [!2345]

